# src/build_data/linux/etc/cron.d/redapid-delete-purged-programs
# /etc/cron.d/redapid-delete-purged-programs
12,42 * * * * root rm -rf /tmp/purged-program-* &> /dev/null
